# deploy/soft-etudes-api.service
[Unit]
Description=Soft Etudes JM API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/soft-etudes/backend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
